# installation/verifiers.conf
# Start verifiers
#

description     "Start R and scipy verifiers server"

start on startup

task
exec python /home/server/scipy-verifier/start_verifiers.py --path=/home/server/scipy-verifier/

